# openclaw/openclaw.service
[Unit]
Description=OpenClaw WebSocket Gateway
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/LPLHackathon-1
EnvironmentFile=/home/ubuntu/LPLHackathon-1/openclaw/.env
Environment="PATH=/home/ubuntu/LPLHackathon-1/venv/bin"
ExecStart=/home/ubuntu/LPLHackathon-1/venv/bin/uvicorn openclaw.openclaw_gateway:app --host 0.0.0.0 --port 18789
Restart=always

[Install]
WantedBy=multi-user.target
